# containers/rootfs/etc/s6-overlay/s6-rc.d/wireproxy/run
#!/command/execlineb -P
with-contenv
s6-notifyoncheck
s6-setuidgid www-data
/usr/bin/wireproxy --silent --config /tmp/wireproxy/wireproxy.conf